# systemd/adsb-logger.service.d/override.conf
[Service]
# Remove the failing recursive chown (it causes the logger to crash-loop)
ExecStartPre=
ExecStartPre=/bin/mkdir -p /opt/adsb-logs
